# templates/cis-report.service.j2
[Unit]
Description=CIS Report Downloader Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=PYTHONPATH={{ charm_dir }}/lib
WorkingDirectory={{ charm_dir }}/src
ExecStart={{ charm_dir }}/venv/bin/python {{ charm_dir }}/src/cis_manager.py --profile-id {{ security_profile_id }}
Restart=on-failure
RestartSec=180
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
